# board/adafruit/adafruit_fruit_jam_rp2350/rootfs_overlay/etc/profile.d/cls.sh
# BusyBox clear applet is not enabled in this image; provide ANSI equivalents.
clear() { printf "\033[2J\033[H"; }
cls() { clear; }
